# runit/godns/run
#!/bin/sh
exec 2>&1
exec godns -c=/path/to/your/godns-dir/config.json