# runit/godns/run
#!/bin/sh
exec 2>&1

[ -r ./conf ] && . ./conf
exec /usr/local/bin/godns -c=/usr/local/etc/godns/config.yaml
