# runit/podman-restart/run
#!/bin/sh
exec 2>&1

[ -r ./conf ] && . ./conf
exec podman --log-level info start --all --filter restart-policy=always
